# game/root/etc/systemd/system/wifi-mode-selector.service
[Unit]
Description=Sélection du mode Wi-Fi (AP ou Client)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-mode-selector.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
